# agent-runtime/agent-runtime.service
[Unit]
Description=AgentShield Agent Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/agentshield/env
ExecStart=/usr/local/bin/agent-runtime
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 安全加固
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/agentshield
PrivateTmp=yes

[Install]
WantedBy=multi-user.target